# src/PiSDWriter/templates/pisdwriter.service.j2
[Unit]
Description=write RaspberryPi SD easiry
After=local-fs.target

[Service]
Type=simple
ExecStart={{ app_execstart_command }}
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
